# build/supervisord.conf
[supervisord]
nodaemon=true
user=root

[unix_http_server]
file=/var/run/supervisor.sock
username=dummy
password=dummy

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock
username=dummy
password=dummy

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:tg]
command=tg 
    --lmdb.mapSize=%(ENV_LMDB_MAPSIZE)s
    --metrics
    --metrics.addr="0.0.0.0"
    --metrics.port="6060"
    --private.api.addr="0.0.0.0:9090"
    %(ENV_EXTRA_OPTS)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
startretries=10

[program:rpcdaemon]
command=rpcdaemon 
    --private.api.addr="localhost:9090"
    --http.addr="0.0.0.0"
    --http.vhosts="*"
    --http.corsdomain="*"
    ps
    --http.api="eth,debug,net"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
startretries=10